# docker/supervisor.d/caddy.ini
[program:caddy]
command=/usr/sbin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/caddy/caddy.stdout.log
stderr_logfile=/var/log/caddy/caddy.stderr.log
