# docker/supervisor.d/php_fpm.ini
[program:php_fpm]
command=/usr/sbin/php-fpm81 -O -R
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/php_fpm/php_fpm.stdout.log
stderr_logfile=/var/log/php_fpm/php_fpm.stderr.log
