# docker/supervisor.d/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/sshd/sshd.stdout.log
stderr_logfile=/var/log/sshd/sshd.stderr.log
