# docker/supervisor.d/webbox_websocket.ini
[program:webbox_websocket]
command=/usr/bin/php /srv/app/app/console --env=prod system:websocket:runserver webbox 8084
user=nobody
group=nobody
numprocs=1
autostart=false
autorestart=true
startsecs=10
startretries=3
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/webbox_websocket/webbox_websocket.stdout.log
stderr_logfile=/var/log/webbox_websocket/webbox_websocket.stderr.log
